# cloudshell.service
[Unit]
Description=cloudshell

[Service]
ExecStart=/usr/bin/cloudshell
[Install]
WantedBy=multi-user.target
